# ops/apps/analyticshub/analyticshub.service.template
[Unit]
Description=@APP_NAME@ backend service
After=network-online.target postgresql.service postgresql-15.service
Wants=network-online.target

[Service]
Type=simple
User=@APP_USER@
Group=@APP_GROUP@
WorkingDirectory=@APP_DIR@
EnvironmentFile=@ENV_FILE@
ExecStart=/usr/bin/java -jar @APP_DIR@/app.jar
SuccessExitStatus=143
Restart=always
RestartSec=5
LimitNOFILE=65536
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
